# debian/gxde-app-upgrader.prerm
#!/bin/bash
systemctl disable spark-update-notifier --now || true