# debian/gxde-app-upgrader.postinst
#!/bin/bash
systemctl daemon-reload
systemctl enable gxde-update-notifier --now || true