# src/usr/lib/systemd/system/gxde-update-notifier.timer
[Unit]  
Description=Timer for GXDE update notifier
  
[Timer]  
OnBootSec=1min
OnUnitActiveSec=1d  
  
[Install]  
WantedBy=timers.target